# supervisord.conf
[supervisord]
nodaemon=true

[program:php-fpm]
command=php-fpm82 -F
autostart=true
autorestart=true

[program:caddy]
command=caddy run --config /etc/caddy/Caddyfile
autostart=true
autorestart=true
